# TPC/Wi-Fi/wifi.service
[Unit]
Description=Start Wi-Fi Drivers
[Service]
ExecStart=/usr/bin/wifi.sh
[Install]
WanterBy=multi-user.target
